# openrc/gnss_share.initd
#!/sbin/openrc-run
# Copyright 2021 Clayton Craft
# SPDX-License-Identifier: GPL-3.0-or-later
command=/usr/bin/gnss_share
command_args="--syslog -v"
command_background=true

depend () {
	provide gpsd
	before geoclue
}
